# sgs-slackbot/crontab
0 * * * * root bash /go/src/github.com/simonmulser/sgs/post/run.sh >> /var/log/cron.log 2>&1
