# sgs-slackbot/crontab
0 * * * * root bash cd /go/sgs-slackbot && ./sgs-slackbot >> /var/log/cron.log 2>&1
  